# systemd/emu-sync.service
[Unit]
Description=Sync ROMs and BIOS files from cloud storage
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/emu-sync sync
Environment=HOME=%h

[Install]
WantedBy=default.target
