# cmd/install_assets/emu-sync.service
[Unit]
Description=Sync ROMs and BIOS files from cloud storage
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=BINARY_PATH sync
Environment=HOME=%h

[Install]
WantedBy=default.target
